# roles/k8s/base/templates/startup.swapoff.service.j2
[Unit]
Description=Swap OFF

[Service]
ExecStart=/usr/sbin/swapoff -a

[Install]
WantedBy=multi.user.target